# root/etc/services.d/plexdrive/run
#!/usr/bin/with-contenv bash

plexdrive_basic_opts="
-c /config/
--cache-file=/config/cache.bolt
--uid=${PUID:-911}
--gid=${PGID:-911}
--umask=022
-o allow_other
"

IFS=" " read -r -a RUN_ARRAY <<< "$RUN_OPTS"

mount_command="plexdrive mount $(echo $plexdrive_basic_opts) /plexdrive ${RUN_ARRAY[@]}"

echo "Executing => $mount_command"
exec s6-setuidgid abc $mount_command
