# root/etc/services.d/ufs/run
#!/usr/bin/with-contenv bash

umask 022

mount_command="ufs"

exec $mount_command
